# deploy/vps/config/systemd/wes-backup.timer
[Unit]
Description=Run WES OS backup daily

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
