# deploy/vps/config/systemd/wes-healthcheck.service
[Unit]
Description=WES OS health check (services + containers + endpoints)
After=wes-infra.service

[Service]
Type=oneshot
ExecStart=/opt/wes/healthcheck.sh
